# pi/armold-watchdog.service
[Unit]
Description=Armold Bridge Watchdog (monitors serial connection and restarts bridge)
After=armold-bridge.service

[Service]
Type=simple
User=pi
Environment=PATH=/opt/ros/jazzy/bin:/home/pi/.local/bin:/usr/bin:/bin
ExecStart=/home/pi/armold_watchdog.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
