# pi/armold-rosbridge.service
[Unit]
Description=Armold ROS 2 WebSocket Bridge (rosbridge_server)
After=armold-bridge.service
Wants=armold-bridge.service

[Service]
Type=simple
User=pi
Environment=PATH=/opt/ros/jazzy/bin:/home/pi/.local/bin:/usr/bin:/bin
ExecStartPre=/bin/sleep 3
ExecStart=/bin/bash -c 'source /opt/ros/jazzy/setup.bash && ros2 launch rosbridge_server rosbridge_websocket_launch.xml'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
